# Tool/src/packaging/debian/xray.conf
# Location /etc/init/xray.conf
# initctl start xray
start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec su -s /bin/bash -c "xray -f /var/log/xray/xray.log" xray
